# deploy/trading-ai-scheduler.service
# systemd unit for the 24/7 agentic scheduler (non-Docker option). Place at
# /etc/systemd/system/trading-ai-scheduler.service then:
#   systemctl enable --now trading-ai-scheduler
[Unit]
Description=Trading-AI Agentic OS scheduler (24/7)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading-ai
Environment=PYTHONIOENCODING=utf-8
Environment=PYTHONPATH=/home/ubuntu/trading-ai
Environment=AOS_PROFILE=micro
ExecStart=/home/ubuntu/trading-ai/.venv/bin/python aos/scheduler.py run-forever
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
